# .config/systemd/user/graphical-environment.target
[Unit]
Description=Graphical environment
BindsTo=graphical-session.target
