# .config/systemd/user/hyprsunset-update.timer
[Unit]
Description=Sunset update for Hyprland

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=hyprsunset-update.service

[Install]
WantedBy=hyprland-session.target
