# .config/systemd/user/tmux-terminal.slice
##
## Tmux (terminal) slice
##

[Service]
CPUWeight=idle
